# systemd/terminal-titles.service
[Unit]
Description=Terminal Titles Updater
After=terminal-manager.service

[Service]
Type=simple
ExecStart=/usr/local/bin/update-terminal-titles.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
