# systemd/workspace-sync.service
[Unit]
Description=Workspace Sync Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/terminal
ExecStart=/usr/bin/node /var/www/terminal/sync-server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
